# conf/crontab
SHELL=/bin/sh
0 0 * * *       /home/qemu-docs/update-docs.sh
